# netrad.init
#!/sbin/runscript

depend() {
        need alsasound localmount
        use net.lo netmount nfsmount esound pulseaudio
}

start() {
    ebegin "Starting netrad"
    env - /usr/local/sbin/netrad ${options} -u ${user} -g ${group} -p ${pidfile} -D "${device}"
    eend $?
}

stop() {
    ebegin "Stopping netrad"
    start-stop-daemon --stop --quiet --pidfile ${pidfile} --user ${user}
    eend $?
}
